# bin/supervisor.ini
[program:gobuild]
command = /home/golang/gobuild/gobuild
directory = /home/golang/gobuild
autostart = true
user = golang
killasgroup = true

stderr_logfile = /var/log/supervisor/gobuild.stdout.log
stdout_logfile = /var/log/supervisor/gobuild.stdout.log
